# crontab
0 * * * * cd /app && /usr/local/bin/python3 /app/app/generate_feed.py >> /var/log/cron.log 2>&1